# upstart.conf
description "Used to start OpenPhotoBooth as a service"
author "Matt Michielsen"
start on runlevel [2]
stop on runlevel [06]
exec /home/dce/src/OpenPhotoBooth/app/start.sh
respawn
